# com.telegram.jobs.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.telegram.jobs.bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/flerov/Desktop/HR ST/my_app/run_bot.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/flerov/Desktop/HR ST/my_app/logs/bot.err.log</string>
    <key>StandardOutPath</key>
    <string>/Users/flerov/Desktop/HR ST/my_app/logs/bot.out.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/flerov/Desktop/HR ST/my_app</string>
</dict>
</plist> 